# scripts/rabbitmq.conf
[program:rabbitmq-server]
command=rabbitmq-server
directory=/srv
autostart=true
autorestart=true
startretries=3
stderr_logfile=/var/log/rabbit-err
stdout-logfile=/var/log/rabbit-out
user=kvorak
